# templates/kubelet.service
[Unit]
Description=kubelet

[Service]
EnvironmentFile=/etc/default/kubelet
ExecStart=/opt/bin/kubelet $KUBELET_OPTS
